# Apps/AIDashApp/Resources/com.tianpli.aidash.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
                       "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tianpli.aidash</string>
    <!--
        SMAppService-bundled LaunchAgents must use `BundleProgram` (a path
        relative to the enclosing .app bundle), not `ProgramArguments` with
        a relative path. With `ProgramArguments=Contents/MacOS/AIDash`,
        launchd resolves the path against its own CWD (`/`), can't find the
        binary, and parks the service in `state = spawn scheduled` forever —
        the mach port is brokered, but the process is never spawned, so
        every CLI connect hangs. See xpc-protocol.md §"Launchd integration".
    -->
    <key>BundleProgram</key>
    <string>Contents/MacOS/AIDash</string>
    <key>MachServices</key>
    <dict>
        <key>com.tianpli.aidash.xpc.v1</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <!--
        `ProcessType=Interactive` tells launchd this agent serves a
        user-facing app (AppKit / SwiftUI menubar) and should not be
        throttled like a background daemon. Required for the menubar UI to
        come up promptly when launchd spawns the agent on-demand.
    -->
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
